# docker/crontab
MAILTO="paton.john@kpmg.nl"
SHELL=/bin/bash

#m  h  dom mon dow command
*   *   *   *   *  echo "I did the thing at $(date)!" >> /usr/src/test.txt
